# ruview_ha/services/bridge/run
#!/usr/bin/with-contenv bashio

export LOG_LEVEL=$(bashio::config 'log_level' 'info' | tr '[:lower:]' '[:upper:]')
export SCAN_INTERVAL=$(bashio::config 'scan_interval' '1')
export INTRUSION_MODE=$(bashio::config 'intrusion_mode' 'false')

bashio::log.info "Starting RuView HA bridge on port ${BRIDGE_PORT}..."
exec python3 -m bridge.main
